# supervisord.conf
[supervisord]
nodaemon=true

[program:hhvm]
command=service hhvm start

[program:nginx]
command=service nginx start